# systemd/coolify-zero.service
[Unit]
Description=Coolify Zero - Zero-downtime deployment manager
Documentation=https://github.com/light-merlin-dark/coolify-zero
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/opt/coolify-zero/bin/coolify-zero.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=coolify-zero

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Environment
Environment="CONFIG_FILE=/etc/coolify-zero/config.yaml"

[Install]
WantedBy=multi-user.target
